# com.auddaily.commodities.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.auddaily.commodities.update</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.12/bin/python3</string>
        <string>/Users/stephanel/AUD-Daily/scripts/Mineral_Commodities_Data_Collection/daily_update.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/stephanel/AUD-Daily</string>
    <key>StandardOutPath</key>
    <string>/Users/stephanel/AUD-Daily/logs/commodities_launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/stephanel/AUD-Daily/logs/commodities_launchd.error.log</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>17</integer>
        <key>Minute</key>
        <integer>0</integer>
        <key>TimeZone</key>
        <string>Australia/Brisbane</string>
    </dict>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
